# backend/newsletter-scanner.timer
[Unit]
Description=AI Newsletter Scanner Timer - Runs every 30 minutes during market hours

[Timer]
# Run at :00 and :30 minutes during market hours (9:30 AM - 4:00 PM ET)
OnCalendar=Mon-Fri 09:30:00
OnCalendar=Mon-Fri 10:00:00
OnCalendar=Mon-Fri 10:30:00
OnCalendar=Mon-Fri 11:00:00
OnCalendar=Mon-Fri 11:30:00
OnCalendar=Mon-Fri 12:00:00
OnCalendar=Mon-Fri 12:30:00
OnCalendar=Mon-Fri 13:00:00
OnCalendar=Mon-Fri 13:30:00
OnCalendar=Mon-Fri 14:00:00
OnCalendar=Mon-Fri 14:30:00
OnCalendar=Mon-Fri 15:00:00
OnCalendar=Mon-Fri 15:30:00
OnCalendar=Mon-Fri 16:00:00
Persistent=true

[Install]
WantedBy=timers.target
