# backend/newsletter-scanner.timer
[Unit]
Description=AI Newsletter Scanner Timer - Runs every 30 minutes during market hours

[Timer]
# Run at minutes 0 and 30 during market hours (9:30 AM - 4:00 PM ET)
OnCalendar=Mon-Fri 09:30..16:00:00,30
OnCalendar=Mon-Fri 09:30..16:00:30,00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
