# services/cycle-scanner.timer
[Unit]
Description=Run Cycle Scanner every 2 hours during market hours
Requires=cycle-scanner.service

[Timer]
# Run at 9:30 AM, 11:30 AM, 1:30 PM, 3:30 PM EST (market hours)
OnCalendar=Mon-Fri *-*-* 09:30:00
OnCalendar=Mon-Fri *-*-* 11:30:00
OnCalendar=Mon-Fri *-*-* 13:30:00
OnCalendar=Mon-Fri *-*-* 15:30:00
Persistent=true

[Install]
WantedBy=timers.target
