# scripts/api-server.service
[Unit]
Description=Options Scanner API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/options-scanner
Environment="PATH=/root/options-scanner/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/options-scanner/venv/bin/python /root/options-scanner/scripts/api_server.py
Restart=always
RestartSec=10
StandardOutput=append:/root/options-scanner/logs/api_server.log
StandardError=append:/root/options-scanner/logs/api_server_error.log

[Install]
WantedBy=multi-user.target
